# artik-plugin/scripts/units/rndis.service
[Unit]
Description=RNDIS USB Function

[Service]
Type=oneshot
ExecStart=/usr/bin/start_rndis.sh

[Install]
WantedBy=multi-user.target
